# ops/systemd/synqra-healthcheck.timer
[Unit]
Description=Run Synqra healthcheck every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=5
Unit=synqra-healthcheck.service

[Install]
WantedBy=timers.target
